# backup.slice
[Unit]
Description=Backup Scripts Slice
Before=slices.target

[Slice]
MemoryAccounting=true
MemoryLimit=1024M
